# crontab
*/1 * * * * /usr/local/Cellar/php@7.2/7.2.23/bin/php /Users/ivosan/oro/bin/console oro:cron --env=prod > /dev/null